# deploy/systemd/solana-api.service
[Unit]
Description=Solana Intel Engine - API Server
Documentation=https://github.com/your-repo/solana-intel-engine
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=solana
Group=solana
WorkingDirectory=/opt/solana-intel-engine

# Environment
EnvironmentFile=/opt/solana-intel-engine/.env

# Start command (Uvicorn with multiple workers)
ExecStart=/opt/solana-intel-engine/.venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000 --workers 2

# Restart policy
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

# Resource limits
MemoryMax=1G
CPUQuota=40%

# Logging
StandardOutput=append:/var/log/solana-intel/api.log
StandardError=append:/var/log/solana-intel/api-error.log

[Install]
WantedBy=multi-user.target
